# scripts/systemd/xon-tranzactions-backend.service
[Unit]
Description=Xon Tranzaksiyalar — Backend (NestJS)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/xon_tranzactions/backend
EnvironmentFile=/var/www/xon_tranzactions/backend/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xon-tranzactions-backend
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
